# vncserver.target
[Unit]
Description=Start VNC server daemons
After=network.target
Wants=network.target

[Install]
WantedBy=multi-user.target
